# config/backerup.service
[Unit]
Description=Backerup backup service

[Service]
ExecStart=/usr/bin/backerup

[Install]
WantedBy=multi-user.target

